# services/browser-workspace-snapshot.service
[Unit]
Description=Snapshot browser workspace layout cache
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/dotfiles/bin/browser-workspace-cache snapshot
