# services/wellbeing-lock.service
[Unit]
Description=Digital Wellbeing Screen Lock
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/home/t/dotfiles/bin/wellbeing-lock
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/t/.Xauthority
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
